# infra/webhook/hermes-webhook.service
[Unit]
Description=Hermes Orchestration GitHub Webhook Receiver
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/hermes-orchestration/webhook
ExecStart=/usr/bin/node /opt/hermes-orchestration/webhook/server.js
Restart=always
RestartSec=5

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/log/hermes-webhook.log /var/log/hermes-orchestration-deploy
ReadOnlyPaths=/etc/eaumalik/hermes-webhook-secret
ProtectHome=true
PrivateTmp=true
MemoryMax=128M

[Install]
WantedBy=multi-user.target
